# sentinel-firstboot.service
[Unit]
Description=Sentinel OS First Boot Hardening
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sentinel-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
